# upstart/ghost.conf
description "ghost blog <blog.theaker.name>"
author "Andrew Theaker <andrew@theaker.name>"

# When to start the service
start on started mongodb and runlevel [2345]

# When to stop the service
stop on shutdown

# Automatically restart process if crashed
respawn
respawn limit 10 5

env name=ghost
env uid=ghost
env gid=ghost
env daemon=/usr/bin/node
env path=/home/ghost/ghost/index.js

script
	export NODE_ENV=production
	exec start-stop-daemon --start --make-pidfile --pidfile /var/run/$name.pid --name $name -c $uid:$gid -x $daemon $path >> /var/log/upstart/$name.log 2>&1
end script
